# .config/systemd/user/restic.service
[Unit]
Description=Restic backup

[Service]
Type=oneshot
ExecStart=restic backup --verbose %h
EnvironmentFile=%h/.config/restic
