# .config/systemd/user/restic.service
[Unit]
Description=Restic backup

[Service]
Type=oneshot
ExecStart=/usr/bin/restic backup --verbose %h
EnvironmentFile=%h/.config/restic
