# .config/systemd/user/tarsnap.service
[Unit]
Description=Tarsnap backup

[Service]
Type=oneshot
ExecStart=%s/.local/bin/tarsnap.sh
