# .config/systemd/user/tarsnap.service
[Unit]
Description=Tarsnap backups

[Service]
Type=oneshot
ExecStart=/home/jarias/.local/bin/tarsnap.sh
